# misc/supervisor_conf.d/newrelic.conf
[program:newrelic]
autostart=false
user=root
priority=2
command=/usr/sbin/nrsysmond -c /etc/newrelic/nrsysmond.cfg -f
