# misc/supervisor_conf.d/sshd.conf
[program:sshd]
priority=1
command=/usr/sbin/sshd -D
